# pladder-irc.service
[Unit]
Description=Pladder Bot IRC Connector

[Service]
Type=notify
ExecStart=/usr/bin/python3 -u -m pladder.irc --systemd --echo localhost 6667 pladder "Pladder IRC Bot" #bot
Restart=on-failure
RestartSec=15min
StartLimitInterval=4h
StartLimitBurst=16
WatchdogSec=30min

[Install]
WantedBy=default.target
