# pladder-upgrade.service
[Unit]
Description=Pladder Bot Upgrader

[Service]
Type=notify
ExecStart=%h/.cache/pladder-venv/bin/gunicorn --access-logfile=- -b localhost:5001 pladder.upgrade:app
SyslogIdentifier=pladder-upgrade

[Install]
WantedBy=default.target
